# luci-app-vnt2/root/etc/init.d/vnt2_client
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=99
STOP=10

. /lib/functions.sh

# ── 复用辅助函数 ──
_opt()  { [ -n "$2" ] && procd_append_param command "$1" "$2"; }
_flag() { [ "$2" = "1" ] && procd_append_param command "$1"; }

_add_input()        { procd_append_param command --input "$1"; }
_add_output()       { procd_append_param command --output "$1"; }
_add_port_mapping() { procd_append_param command --port-mapping "$1"; }

start_service() {
	config_load vnt2

	local enabled; config_get_bool enabled client enabled 0
	[ "$enabled" = 1 ] || return 0

	local cli_bin web_bin
	config_get cli_bin globals client_bin '/usr/bin/vnt2_cli'
	config_get web_bin globals web_bin    '/usr/bin/vnt2_web'
	[ -x "$cli_bin" ] || { logger -t vnt2_client "$cli_bin 未找到"; return 1; }

	# 必填项
	local server network_code ip
	config_get server       client server
	config_get network_code client network_code
	config_get ip           client ip

	[ -z "$server" ]       && { logger -t vnt2_client "缺少服务器地址"; return 1; }
	[ -z "$network_code" ] && { logger -t vnt2_client "缺少网络编号"; return 1; }
	[ -z "$ip" ]           && { logger -t vnt2_client "缺少虚拟IP"; return 1; }

	# 可选项
	local password device_name device_id cert_mode tun_name mtu ctrl_port
	local rtx compress fec no_punch no_nat no_tun allow_mapping
	config_get password    client password
	config_get device_name client device_name
	config_get device_id   client device_id
	config_get cert_mode   client cert_mode
	config_get tun_name    client tun_name
	config_get mtu         client mtu
	config_get ctrl_port   client ctrl_port
	config_get_bool rtx          client rtx 0
	config_get_bool compress     client compress 0
	config_get_bool fec          client fec 0
	config_get_bool no_punch     client no_punch 0
	config_get_bool no_nat       client no_nat 0
	config_get_bool no_tun       client no_tun 0
	config_get_bool allow_mapping client allow_mapping 0

	# 设备名称：未设置时使用主机名
	[ -z "$device_name" ] && device_name=$(cat /proc/sys/kernel/hostname 2>/dev/null)

	procd_open_instance vnt2_client
	procd_set_param command "$cli_bin"

	# 必填参数
	procd_append_param command --server "$server"
	procd_append_param command --network-code "$network_code"
	procd_append_param command --ip "$ip"

	# 可选参数
	_opt  --password    "$password"
	_flag --rtx         "$rtx"
	_flag --compress    "$compress"
	_flag --fec         "$fec"

	config_list_foreach client input  _add_input
	config_list_foreach client output _add_output

	_opt  --device-name   "$device_name"
	_opt  --device-id     "$device_id"
	_flag --no-punch      "$no_punch"
	_opt  --cert-mode     "$cert_mode"
	_opt  --tun-name      "$tun_name"
	_flag --no-nat        "$no_nat"
	_flag --no-tun        "$no_tun"
	_flag --allow-mapping "$allow_mapping"

	config_list_foreach client port_mapping _add_port_mapping

	_opt --mtu       "$mtu"
	_opt --ctrl-port "$ctrl_port"

	procd_set_param respawn
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance

	# ── vnt2_web 可选实例 ──
	local web_enabled web_addr web_conf
	config_get_bool web_enabled client web_enabled 0
	if [ "$web_enabled" = 1 ] && [ -x "$web_bin" ]; then
		config_get web_addr client web_addr '0.0.0.0:29870'
		config_get web_conf client web_conf ''

		procd_open_instance vnt2_web
		procd_set_param command "$web_bin" --addr "$web_addr"
		[ -n "$web_conf" ] && procd_append_param command --conf "$web_conf"
		procd_set_param respawn
		procd_set_param stdout 1
		procd_set_param stderr 1
		procd_close_instance
	fi
}

service_triggers() { procd_add_reload_trigger vnt2; }
